# templates/systemd.service.j2
[Unit]
Description=Perform {{ owncloud_vendor }} tasks

[Service]
Type=oneshot
User=apache
Group=apache
ExecStart=/usr/bin/php -f {{ owncloud_destination }}/cron.php
