# templates/systemd.timer.j2
[Unit]
Description=Perform {{ owncloud_vendor }} tasks

[Timer]
OnCalendar=*:0/15
